# alienware-audio-fix.service
[Unit]
Description=Alienware Audio Fix (Topology & Driver Reload)
After=sound.target
Wants=sound.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/alienware-audio-boot-fix.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
